# crossbackup.timer
[Unit]
Description=Schedule crossbackup service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
